# deploy/boomerang.service
[Unit]
Description=Boomerang backup appliance
After=network.target

[Service]
Type=simple
User=boomerang
Group=boomerang
Environment=BOOMERANG_DATA_DIR=/var/lib/boomerang
Environment=BOOMERANG_LISTEN=0.0.0.0:8080
ExecStart=/usr/local/bin/boomerang
Restart=on-failure
RestartSec=3
# NoNewPrivileges must stay off so the boomerang user can run the limited
# passwordless sudo rule for /usr/local/sbin/boomerang-update (in-app updates).
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/boomerang
PrivateTmp=true

[Install]
WantedBy=multi-user.target
